# packages/ubuntu/debian/postinst
#!/bin/sh -e

case "$1" in
    configure)
	# continue below
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
	exit 0
    ;;

    *)
	echo "postinst called with unknown argument \`$1'" >&2
	exit 0
    ;;
esac

umask 022

ntopng_found=`/bin/grep ^ntopng /etc/passwd|wc -l`
[ "$ntopng_found" -eq "0" ] && /usr/sbin/useradd ntopng -p ntopng

echo "Rebuilding ld cache..."
/sbin/ldconfig

echo "Adding the ntopng startup script"
update-rc.d ntopng defaults 93 >/dev/null

echo "Making the /etc/ntopng directory..."

mkdir -p /etc/ntopng/


exit 0
